# auction-closer/crontab
# Based on Kulatunga (2021)
SHELL=/bin/bash
BASH_ENV=/etc/environment

 * * * * * cd /usr/src/app && npm start > /proc/1/fd/1 2>/proc/1/fd/2
